# launchd/com.nanoclaw.kanban-renderer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nanoclaw.kanban-renderer</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{PROJECT_ROOT}}/scripts/kanban-renderer.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{PROJECT_ROOT}}</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>25</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PATH</key>
    <string>{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>KANBAN_RENDER_GROUP</key>
    <string>main</string>
  </dict>

  <key>StandardOutPath</key>
  <string>{{PROJECT_ROOT}}/logs/kanban-renderer.log</string>
  <key>StandardErrorPath</key>
  <string>{{PROJECT_ROOT}}/logs/kanban-renderer.error.log</string>
</dict>
</plist>
