# launchd/com.nanoclaw.stuck-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nanoclaw.stuck-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{NODE_PATH}}</string>
    <string>{{PROJECT_ROOT}}/scripts/stuck-monitor.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{PROJECT_ROOT}}</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PATH</key>
    <string>{{HOME}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>STUCK_MONITOR_GROUP</key>
    <string>main</string>
  </dict>

  <key>StandardOutPath</key>
  <string>{{PROJECT_ROOT}}/logs/stuck-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>{{PROJECT_ROOT}}/logs/stuck-monitor.error.log</string>
</dict>
</plist>
